# scripts/backup/clean_backups.timer
[Unit]
Description=Clean up old backups

[Timer]
OnCalendar=monthly
Persistent=true

[Install]
WantedBy=timers.target
